# build/supervisord.conf
[supervisord]
nodaemon=true

[program:web-server]
command=/usr/bin/python2.7 /opendc/opendc-web-server/main.py /opendc/keys.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:simulator]
directory=/opendc/opendc-simulator/opendc-integration-jpa
command=/opendc/opendc-simulator/gradlew run
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

