# build/supervisord.conf
[supervisord]
nodaemon=true

[program:web-server]
command=/usr/bin/python2.7 /opendc/opendc-web-server/main.py /opendc/keys.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:simulator]
command=/bin/sh -c '/opendc/opendc-simulator/Simulator/simulator /data/database/opendc.db'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
